# etc/cron.d/letterboxd-netflix
30 20 * * * letterboxd /opt/letterboxd/sbin/netflix-update
