# etc/cron.d/letterboxd-availability
# Cron jobs for the primary app server
# To update availability services

# Netflix - daily
30 20 * * * letterboxd /opt/letterboxd/sbin/netflix-update

# iTunes - Sunday 0:00, which appears to be a number of hours after the weekly data update
0 0 * * 0 letterboxd /opt/letterboxd/sbin/itunes-download-update

# Amazon - Wednesday 0:00
0 0 * * 3 letterboxd /opt/letterboxd/sbin/amazon-update
