# etc/cron.d/letterboxd-db
# Cron jobs for the primary database server

30 22 * * * letterboxd /opt/letterboxd/libexec/postgresql/maintain-database.sh
