# etc/cron.d/pganalyze
*/10 * * * * pganalyze /usr/bin/pganalyze-collector --cron
